# environment/school-market.crontab
basePath=/home/gorod
0 */12 * * * cd $basePath && ./console/cron/updateSearch.sh
0 */12 * * * cd $basePath && ./console/cron/updateRank.sh
